# files/system/usr/lib/systemd/system/cec-tv-on.service
[Unit]
Description=Wake the TV and switch it to this PC over HDMI-CEC
After=cec-follower.service
Wants=cec-follower.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/plasma-machine/cec-tv-on

[Install]
WantedBy=multi-user.target
